# deploy/vantrufree-leads.service
[Unit]
Description=VANTRUFREE Leads API
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/vantrufree/backend
ExecStart=/usr/bin/node --env-file=.env server.mjs
Restart=always
RestartSec=3
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
